# deploy/supervisord.conf
; Run the 3 services side-by-side in one container (all on localhost).
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

; Headless Chrome with CDP — the SGE scraper connects here (localhost:9222).
; The SGE login cookie lives in /profile (mount it; see DEPLOY.md).
[program:chrome]
command=/usr/bin/chromium --headless=new --no-sandbox --disable-gpu --disable-dev-shm-usage --remote-debugging-port=9222 --user-data-dir=/profile --no-first-run --no-default-browser-check
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Python scraper sidecar — TikTok + Instagram (localhost:8900).
[program:scraper]
command=/opt/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8900
directory=/app/scraper-service
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; The Next app — the platform (0.0.0.0:3000, reachable from outside).
[program:next]
command=npm run start
directory=/app
autorestart=true
priority=30
environment=PORT="3000",HOSTNAME="0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
